# pipeline/systemd/review-detector.timer
[Unit]
Description=Run the review-detector every 5 minutes
Documentation=file:/home/cvps/documentation/hermes-phase-g-plan.md

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
