# systemd/clawhip.service.d/10-gjc-relay.conf
# clawhip drop-in (installed as /etc/systemd/system/clawhip.service.d/10-gjc-relay.conf).
# Ordering only — the base clawhip.service unit is NOT edited. Ensures the relay is
# up before clawhip starts sending, and is pulled in as a dependency.
[Unit]
After=gjc-relay.service
Wants=gjc-relay.service
