# systemd/clawhip.service
[Unit]
Description=clawhip daemon (event router -> Discord)
# User-scope unit: no network-online.target exists in the user manager; the
# daemon's Restart loop self-heals pre-network starts (RestartSec=5).

[Service]
Type=simple
WorkingDirectory=%h
Environment=RUST_LOG=info
Environment=CLAWHIP_CONFIG=%h/.clawhip/config.toml
EnvironmentFile=%h/.clawhip/clawhip.env
ExecStart=%h/.cargo/bin/clawhip start
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
