# systemd/gjc-relay-heartbeat.service
[Unit]
Description=gjc-relay token-cache heartbeat (self-priming no-op inbound via clawhip)
Documentation=file:{{FLEET_REPO}}/docs/35-gjc-relay.md
After=clawhip.service gjc-relay.service

[Service]
Type=oneshot
EnvironmentFile=%h/.gjc-relay/relay.env
ExecStart=%h/.gjc-relay/relay-heartbeat.sh
# Lightweight hardening only, matching gjc-relay.service: this host restricts
# unprivileged user namespaces via AppArmor, so namespace-based sandboxing
# (ProtectSystem/ProtectHome/PrivateTmp) is deliberately ABSENT here too — see
# gjc-relay.service for the full rationale. The directives below are all
# honored without namespaces.
NoNewPrivileges=true
RestrictRealtime=true
LockPersonality=true
SystemCallArchitectures=native
RestrictNamespaces=true
MemoryDenyWriteExecute=true
