# systemd/gjc-worktree-janitor.timer
[Unit]
Description=Run the gjc worktree janitor every 2 minutes (crash-net)
Documentation=file:{{FLEET_REPO}}/docs/40-gjc-bot-automation.md

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=15s
Persistent=false

[Install]
WantedBy=timers.target
