# systemd/gjc-dlq-watch.service
[Unit]
Description=gjc DLQ-bury watcher (out-of-band alert on clawhip silent notification loss)
Documentation=file:{{FLEET_REPO}}/docs/35-gjc-relay.md
After=clawhip.service

[Service]
Type=simple
EnvironmentFile=-%h/.gjc-relay/relay.env
ExecStart=%h/.gjc-relay/dlq-watch.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
