# systemd/gjc-relay-alert.service
[Unit]
Description=gjc-relay out-of-band failure alarm (direct-to-Discord; bypasses clawhip + relay)
Documentation=file:{{FLEET_REPO}}/docs/35-gjc-relay.md

[Service]
Type=oneshot
EnvironmentFile=-%h/.gjc-relay/relay.env
ExecStart=%h/.gjc-relay/alert.sh
