# systemd/gjc-relay.service
[Unit]
Description=gjc-relay (clawhip -> Discord transform relay, loopback 127.0.0.1)
Documentation=file:{{FLEET_REPO}}/docs/35-gjc-relay.md
# Critical-path supervision: never give up restarting. A permanently-down relay
# silently DLQ-buries every clawhip notification (clawhip does not retry transport
# errors and its DLQ is a non-persistent in-memory Vec).
StartLimitIntervalSec=0
# Best-effort out-of-band alarm. NOTE: with Restart=always + StartLimitIntervalSec=0
# the unit is restarted on every exit and can never reach the `failed` state, so this
# rarely fires by design. The OPERATIVE alarm for real notification loss is
# gjc-dlq-watch.service, which watches clawhip's journal for "clawhip dlq bury:" and
# alerts direct-to-Discord (bypassing relay + clawhip). We keep infinite restart on
# purpose — a finite start-limit would resurrect the "permanently down, silent loss" risk.
OnFailure=gjc-relay-alert.service

[Service]
Type=simple
WorkingDirectory=%h/.gjc-relay
EnvironmentFile=%h/.gjc-relay/relay.env
ExecStart=%h/.gjc-relay/gjc-relay
Restart=always
RestartSec=1
# Lightweight hardening (relay only reads config + makes outbound TLS, writes
# nothing to disk). Namespace-based sandboxing (ProtectSystem/ProtectHome/
# PrivateTmp) is deliberately ABSENT: in a user manager those need unprivileged
# user namespaces, which Ubuntu >= 24.04 restricts via AppArmor
# (kernel.apparmor_restrict_unprivileged_userns=1) — a hard start failure on the
# fleet's single point of failure is worse than the lost sandbox. The directives
# below are all honored without namespaces.
NoNewPrivileges=true
RestrictRealtime=true
LockPersonality=true
SystemCallArchitectures=native
RestrictNamespaces=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=default.target
