# systemd/issue-spool-adapter.path
[Unit]
Description=Watch the clawhip issue spool and run the intake adapter
Documentation=file:{{FLEET_REPO}}/docs/40-gjc-bot-automation.md

[Path]
# clawhip's localfile sink appends github.issue-opened records here; each append
# triggers the adapter. The adapter re-reads the whole spool and dedups via the
# ledger, so a coalesced/missed trigger is self-correcting.
PathModified=%h/.gjc-bot/issue-spool.jsonl
Unit=issue-spool-adapter.service

[Install]
WantedBy=paths.target
