# systemd/merge-gate.service
[Unit]
Description=Advisory merge gate — post MERGE_READY/REQUEST_CHANGES on CI-green bot PRs (comment-only)
Documentation=file:{{FLEET_REPO}}/docs/40-gjc-bot-automation.md

[Service]
Type=oneshot
Environment=PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-%h/.gjc-bot/gjc-bot.env
ExecStart={{FLEET_REPO}}/pipeline/review/merge-gate.sh
Nice=15
