# systemd/review-detector.service
[Unit]
Description=Poll bot PRs for augmentcode[bot] reviews-with-suggestions -> launch AI Code Review Handler
Documentation=file:{{FLEET_REPO}}/docs/40-gjc-bot-automation.md

[Service]
Type=oneshot
Environment=PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-%h/.gjc-bot/gjc-bot.env
ExecStart={{FLEET_REPO}}/pipeline/review/review-detector.sh
Nice=10
# The detector fire-and-forgets a detached handler run (review-run.sh -> setsid).
# KillMode=process keeps that run alive after this oneshot's main process exits
# (default control-group would reap the whole cgroup).
KillMode=process
